# roles/kubernetes/postinstall/files/renew-certs-kubeadm.timer
[Unit]
Description=Timer of Certificates rotation service for Kubernetes
[Timer]
OnCalendar=quarterly
Persistent=1
