# roles/network_plugin/fan/templates/cni-dhcp.service.j2
[Service]
ExecStart=/opt/cni/bin/dhcp daemon
