# roles/network_plugin/fan/templates/cni-dhcp.socket.j2
[Socket]
ListenStream=/run/cni/dhcp.sock

[Install]
WantedBy=sockets.target
