# roles/kubernetes/control_plane/templates/k8s-certs-renew.service.j2
[Unit]
Description=Renew K8S control plane certificates

[Service]
Type=oneshot
ExecStart={{ bin_dir }}/k8s-certs-renew.sh
