# roles/kubernetes/control_plane/templates/k8s-certs-renew.timer.j2
[Unit]
Description=Timer to renew K8S control plane certificates

[Timer]
OnCalendar={{ auto_renew_certificates_systemd_calendar }}

[Install]
WantedBy=multi-user.target
