# roles/kubernetes/control-plane/templates/k8s-certs-renew.timer.j2
[Unit]
Description=Timer to renew K8S control plane certificates

[Timer]
# First Monday of each month
OnCalendar=Mon *-*-1..7 03:{{ groups['kube_control_plane'].index(inventory_hostname) }}0:00

[Install]
WantedBy=multi-user.target
