# roles/kubernetes/node/templates/kubelet.host.service.j2
[Unit]
Description=Kubernetes Kubelet Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=docker.service
Wants=docker.socket

[Service]
EnvironmentFile=-{{kube_config_dir}}/kubelet.env
{% if ansible_os_family in ["CoreOS", "Container Linux by CoreOS"] %}
Environment="PATH=/usr/local/sbin:/usr/sbin:/usr/bin:/sbin:/bin:{{ bin_dir }}"
{% endif %}
ExecStart={{ bin_dir }}/kubelet \
		$KUBE_LOGTOSTDERR \
		$KUBE_LOG_LEVEL \
		$KUBELET_API_SERVER \
		$KUBELET_ADDRESS \
		$KUBELET_PORT \
		$KUBELET_HOSTNAME \
		$KUBE_ALLOW_PRIV \
		$KUBELET_ARGS \
		$DOCKER_SOCKET \
		$KUBELET_NETWORK_PLUGIN \
		$KUBELET_CLOUDPROVIDER
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
