# roles/kubernetes/postinstall/templates/renew-certs-kubeadm.timer
[Unit]
Description=Timer of Certificates rotation service for Kubernetes
[Timer]
OnCalendar={{ postinstall_rotation_period }}
Persistent=1
